# scripts/com.nikkei-trade-journal.fetch.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd agent: runs scripts/fetch-market-data.sh every 15 min.
     Install:
       1) Replace ABSOLUTE_REPO_PATH and YOUR_HOME below.
       2) cp this file to ~/Library/LaunchAgents/com.nikkei-trade-journal.fetch.plist
       3) launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.nikkei-trade-journal.fetch.plist
     Requires the backend (make server) to be running on BASE_URL. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nikkei-trade-journal.fetch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>ABSOLUTE_REPO_PATH/scripts/fetch-market-data.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>BASE_URL</key>
    <string>http://localhost:8080</string>
    <key>INTERNAL_SECRET</key>
    <string>dev-secret</string>
    <key>SYMBOLS</key>
    <string>^N225 NKD=F</string>
  </dict>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>YOUR_HOME/Library/Logs/nikkei-trade-journal-fetch.log</string>
  <key>StandardErrorPath</key>
  <string>YOUR_HOME/Library/Logs/nikkei-trade-journal-fetch.log</string>
</dict>
</plist>
